# debian/cron.d
MAILTO="peter.howe@emc.com"

*/5 * * * * root [ -x /usr/sbin/repoSync ] && /usr/sbin/repoSync

