# debian/cron.d
MAILTO="peter.howe@emc.com"

*/15 * * * * root [ -x /usr/local/packrat/cron/repoSync ] && /usr/local/packrat/cron/repoSync

